# terraform/stacks/prod/tee/ansible/playbooks/templates/tee-api.service.j2
[Unit]
Description=Tearleads Tee API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ tee_api_user }}
Group={{ tee_api_group }}
WorkingDirectory={{ tee_api_install_dir }}
Environment=NODE_ENV=production
EnvironmentFile={{ tee_api_config_dir }}/tee-api.env
ExecStart=/usr/bin/node {{ tee_api_install_dir }}/server.cjs
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ tee_api_data_dir }}

[Install]
WantedBy=multi-user.target
