# ansible/playbooks/templates/auto-patch.timer.j2
[Unit]
Description=Daily automatic patching timer

[Timer]
OnCalendar=*-*-* {{ patch_time }}:00 {{ patch_timezone }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
