# ansible/playbooks/templates/sync-last-active.timer.j2
[Unit]
Description=Daily sync of lastActiveAt from Redis to PostgreSQL

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
