# ansible/playbooks/templates/tearleads-api.service.j2
[Unit]
Description=Tearleads API Server
Documentation=https://github.com/a2f0/tearleads
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tearleads-api
ExecStart=/usr/bin/node dist/server.cjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tearleads-api

# Environment
Environment=NODE_ENV=production
Environment=PORT=5001
EnvironmentFile=/opt/tearleads-api/.env

# COMPLIANCE_SENTINEL: TL-SVC-001 | control=service-sandboxing
# Security hardening - systemd sandbox
# Framework mappings:
# - SOC2: CC6.1 (Logical Access Security)
# - NIST SP 800-53: SC-7, AC-6 (Boundary Protection, Least Privilege)
# - HIPAA: 164.312(a)(1) (Access Control)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tearleads-api
PrivateTmp=true

# Additional security hardening
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true

# Restrict network address families (IPv4/IPv6 only)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Limit system call exposure
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @swap @reboot @raw-io

# UMask for file creation
UMask=0077

[Install]
WantedBy=multi-user.target
