# ansible/playbooks/templates/tearleads-api.service.j2
[Unit]
Description=Tearleads API Server
Documentation=https://github.com/a2f0/tearleads
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tearleads-api
ExecStart=/usr/bin/node dist/server.cjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tearleads-api

# Environment
Environment=NODE_ENV=production
Environment=PORT=5001
EnvironmentFile=/opt/tearleads-api/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tearleads-api
PrivateTmp=true

[Install]
WantedBy=multi-user.target
